# acl-watchd.conf
# /etc/event.d/acl-watchd
#
# Upstart job for the ACL watch daemon.
# Starts after the filesystem is writable; restarts if it crashes.
#
# Install:
#   novacom put file:///etc/event.d/acl-watchd < acl-watchd.conf
#   novacom run -- file://sbin/start acl-watchd

description "ACL Watch Daemon"
author      "ACL Manager"

# Start once the system is up enough to have /media/internal mounted.
# "stopped rcS" fires after the early boot scripts finish.
start on stopped rcS

respawn

exec /usr/bin/acl-watchd
